# launchd/com.local.mailcleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.local.mailcleanup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/mymac/.pyenv/shims/python3</string>
        <string>/Users/mymac/Scripts/mail-imap/cleanup.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mymac/Scripts/mail-imap</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/mymac/Scripts/mail-imap/logs/cleanup.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mymac/Scripts/mail-imap/logs/cleanup.err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
